# dropboxd.service
[Unit]
Description=Dropbox Service

[Service]
ExecStart=/home/saminiir/.dropbox-dist/dropboxd

[Install]
WantedBy=multi-user.target
